# supervisord-default.conf
[program:apache]
user=www-data
command=/usr/sbin/httpd -DFOREGROUND
numprocs=1
autostart=true
autorestart=true

[program:cron]
user=www-data
command=cron -f -L 15
numprocs=1
autostart=true
autorestart=true
