# supervisord-apache2.conf
[program:apache2]
user=www-data
command=/usr/sbin/apache2 -DFOREGROUND
numprocs=1
autostart=true
autorestart=true
