# rootfs/etc/supervisor.d/httpd.ini
[program:httpd]
command=/usr/local/bin/httpd-foreground 
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisord/httpd-error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisord/httpd-access.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
